# dot_config/systemd/user/serena.service.tmpl
[Unit]
Description=Serena MCP Server
After=network.target

[Service]
Type=simple
User=sboynton
WorkingDirectory={{ .chezmoi.homeDir }}
ExecStart={{ .chezmoi.homeDir }}/.local/bin/uvx --from git+https://github.com/oraios/serena serena start-mcp-server --transport streamable-http --port 9121 --context ide-assistant
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
