# dot_config/systemd/user/waybar.service
[Unit]
Description=Waybar

[Service]
Type=simple
ExecStart=/usr/bin/waybar
Restart=always
RestartSec=10s

[Install]
WantedBy=default.target
